# spring/spring-boot-web/src/main/resources/META-INF/boot-web.service
[Unit]
Description=Spring Boot Web Application

[Service]
ExecStart=/usr/bin/java -jar location_of_jar_file.jar --spring.config.location=location_of_config.properties --spring.profiles.active=profile
User=jelly

[Install]
WantedBy=multi-user.target